# scripts/launchd/com.cucuradu.trading-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cucuradu.trading-dispatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/raducucu/Documents/Python/Trading Bot/scripts/cron_dispatch.sh</string>
  </array>

  <!-- Run every 5 minutes; the dispatcher itself checks ET wall-clock and
       fires the matching routine. DST-correct on both US and EU sides. -->
  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <false/>

  <!-- launchd has a stripped PATH; export everything claude needs. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/raducucu/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/raducucu</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/raducucu/Documents/Python/Trading Bot</string>

  <key>StandardOutPath</key>
  <string>/Users/raducucu/Documents/Python/Trading Bot/cache/dispatch/launchd.out</string>

  <key>StandardErrorPath</key>
  <string>/Users/raducucu/Documents/Python/Trading Bot/cache/dispatch/launchd.err</string>
</dict>
</plist>
